# dell-fan-control-service/DEBIAN/prerm
#!/bin/sh
systemctl stop dell-fan-control-service.timer
systemctl disable dell-fan-control-service.timer