# dell-fan-control-service/lib/systemd/system/dell-fan-control-service.timer
[Unit]
Description=Start dell-fan-control every minute

[Timer]
OnBootSec=60s
OnUnitActiveSec=59s
AccuracySec=2s

[Install]
WantedBy=timers.target
