# dell-fan-control-service/lib/systemd/system/dell-fan-control-service.service
[Unit]
Description=Dell Fan Control

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/bin/dell-fan-control.sh